# support/admin/server
#!/sbin/openrc-run

name="serve repos"
command="/usr/bin/python3.8"
command_args="/root/newhope/cli pkg serve repo2 --fr /media/storage/pkg_repo"
command_background=true
pidfile="/run/server.pid"
