# support/admin/syncer
#!/sbin/openrc-run

name="sync build ad prod repos"
command="/usr/bin/python3.8"
command_args="/root/newhope/cli pkg sync repo --fr /media/build/dev/r --to /media/storage/pkg_repo"
command_background=true
pidfile="/run/syncer.pid"
