# support/admin/runit/fwd_22/run
#!/bin/sh

ssh -N -R *2200:localhost:22 root@138.68.80.104
