# support/admin/runit/scheduler/run
#!/bin/sh

/home/pg83/newhope/cli cmd scheduler
