# support/admin/runit/scheduler/run
#!/bin/sh

exec /home/pg83/newhope/cli cmd scheduler >& /media/storage/sched.log
