# support/admin/runit/fwd_80/run
#!/bin/sh

ssh -N -R 80:localhost:80 root@138.68.80.104
