# support/admin/runit/fwd_80/run
#!/bin/sh

exec unshare --fork --pid --kill-child autossh -M 0 -N -R 80:localhost:80 root@138.68.80.104
